# main-agent/openrc/agent
#!/sbin/openrc-run

command="/app/main-agent"
command_args=""
supervisor="supervise-daemon"
name="Grobuzin agent"

depend() {
    need net
    after firewall
}

stop() {
    ebegin "Stopping Main Agent"
    start-stop-daemon --stop --exec $command
    eend $?
}
