# main-agent/openrc/agent
#!/sbin/openrc-run

command="/app/main-agent"
command_args=""
pidfile="/var/run/main-agent.pid"
name="Grobuzin agent"

depend() {
    need net
    after firewall
}

start_pre() {
    # Commands to run before starting the service
}

start() {
    ebegin "Starting Main Agent"
    start-stop-daemon --start --exec $command -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping Main Agent"
    start-stop-daemon --stop --exec $command
    eend $?
}

status() {
    status -p $pidfile $command
}